# config/service/php55/run
#!/bin/bash

exec /usr/sbin/php-fpm5.5 -F > /dev/null 2>&1 & wait